# config/clawbox-browser.service
[Unit]
Description=ClawBox Browser (Chromium with CDP)
Requires=clawbox-vnc.service
After=clawbox-vnc.service

[Service]
Type=simple
Environment=DISPLAY=:99
Environment=CDP_PORT=18800
Environment=HOME=/home/clawbox
ExecStart=/home/clawbox/clawbox/scripts/launch-browser.sh
Restart=no
TimeoutStartSec=30
# KillMode=control-group already terminates all child processes (chrome_crashpad,
# the main chrome, dbus-launch helpers) when the unit stops — no ExecStopPost
# pgrep needed (and the prior pattern was broad enough to kill unrelated procs).
KillMode=control-group

[Install]
WantedBy=multi-user.target
