# config/clawbox-ap.service
[Unit]
Description=ClawBox WiFi Access Point and Captive Portal
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Root is required for iptables NAT rules and NetworkManager AP management
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE
Environment=CLAWBOX_HOME=/home/clawbox/clawbox
ExecStartPre=/usr/bin/test -x ${CLAWBOX_HOME}/scripts/start-ap.sh
ExecStartPre=/usr/bin/test -x ${CLAWBOX_HOME}/scripts/stop-ap.sh
ExecStart=${CLAWBOX_HOME}/scripts/start-ap.sh
ExecStop=${CLAWBOX_HOME}/scripts/stop-ap.sh

[Install]
WantedBy=multi-user.target
