# config/clawbox-performance.service
[Unit]
Description=ClawBox MAXN Power Mode and Jetson Clocks
After=multi-user.target nvpmodel.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c 'MAXN_ID=$(grep -oP "POWER_MODEL ID=\\K\\d+(?=\\s+NAME=MAXN)" /etc/nvpmodel.conf | tail -1); nvpmodel -m ${MAXN_ID:-0} && jetson_clocks'
ExecStop=/bin/bash -c 'jetson_clocks --restore 2>/dev/null || true'

[Install]
WantedBy=multi-user.target
