# config/clawbox-gateway.service
[Unit]
Description=ClawBox OpenClaw Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=clawbox
WorkingDirectory=/home/clawbox
ExecStartPre=/home/clawbox/clawbox/scripts/gateway-pre-start.sh
ExecStart=/home/clawbox/.npm-global/bin/openclaw gateway --allow-unconfigured --bind loopback
Restart=always
RestartSec=5
Environment=HOME=/home/clawbox
Environment=NODE_ENV=production
Environment=NODE_OPTIONS=--max-old-space-size=1024
Environment=OPENAI_TTS_BASE_URL=http://localhost:8880/v1
Environment=OPENCLAW_VERBOSE=1
# Route ClawHub traffic through our cached edge proxy by default.
# Override with `systemctl edit clawbox-gateway.service` or by setting
# OPENCLAW_CLAWHUB_URL in /home/clawbox/clawbox/data/network.env (loaded next).
Environment=OPENCLAW_CLAWHUB_URL=https://openclawhardware.dev
EnvironmentFile=-/home/clawbox/clawbox/data/network.env

[Install]
WantedBy=multi-user.target
