# config/clawbox-root-update@.service
[Unit]
Description=ClawBox Root Update Step (%i)

[Service]
Type=oneshot
ExecStart=/home/clawbox/clawbox/scripts/root-update-step.sh %i
TimeoutStartSec=600
